# whatsapp-relay/relay-monitor.service
[Unit]
Description=ProductFlow Relay (Baileys)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/pi/ProductFlow/whatsapp-relay
ExecStart=/usr/bin/node relay-baileys.js
Restart=on-failure
RestartSec=30
EnvironmentFile=/home/pi/ProductFlow/whatsapp-relay/.env

[Install]
WantedBy=multi-user.target
